# pagl/rootfs/etc/s6-overlay/s6-rc.d/svc-loki/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Loki
# ==============================================================================

bashio::log.info "Starting Loki..."

exec loki \
  -config.file=/etc/loki/loki.yml
