# librespot/rootfs/etc/services.d/spotifyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs librespot
# ==============================================================================
bashio::log.info 'Starting librespot...'

# Run librespot
exec /usr/bin/librespot
