# librespot/rootfs/etc/services.d/spotifyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs librespot
# ==============================================================================

name=$(bashio::config 'name')
bitrate=$(bashio::config 'bitrate')
username=$(bashio::config 'username')
password=$(bashio::config 'password')
volume=$(bashio::config 'volume')

# Create snapfifo dir
bashio::log.info 'Starting librespot...'

mkdir -p /share/snapfifo

# Run librespot
bashio::log.info 'Starting librespot...'

exec /usr/local/bin/librespot  --name ${name} --backend pipe --device /share/snapfifo/librespot --bitrate ${bitrate} --disable-audio-cache --enable-volume-normalisation --linear-volume --initial-volume=${volume} --username ${username} --password ${password} --zeroconf-port 35535
