# qbittorrent/rootfs/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the qbittorrent daemon
# ==============================================================================

bashio::log.info "Starting qbittorrent..."
exec /usr/bin/qb_up
