# pagl/rootfs/etc/s6-overlay/s6-rc.d/svc-grafana/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Grafana
# ==============================================================================

bashio::log.info "Starting Grafana..."

exec grafana server \
  --config=/etc/grafana/grafana.ini \
  --homepath=/usr/share/grafana
