# immich/rootfs/etc/s6-overlay/s6-rc.d/svc-server/finish
#!/usr/bin/with-contenv bash
# Forcefully clean up orphaned Immich processes and wait for port release

port="${SERVER_PORT:-2283}"
max_attempts=20
attempt=0

echo "svc-server/finish: killing Immich processes on port ${port}" >&2

# Kill all node processes
killall -9 node 2>/dev/null || true
killall -9 start.sh 2>/dev/null || true
sleep 1

# Verify port is released, with retries
while [ $attempt -lt $max_attempts ]; do
    if ! lsof -i :${port} 2>/dev/null | grep -q LISTEN; then
        echo "svc-server/finish: port ${port} is now free" >&2
        exit 0
    fi
    
    echo "svc-server/finish: port ${port} still in use, waiting... (${attempt}/${max_attempts})" >&2
    sleep 1
    ((attempt++))
done

echo "svc-server/finish: WARNING timeout waiting for port ${port} to release" >&2
# Kill anything still listening on that port
lsof -ti :${port} 2>/dev/null | xargs kill -9 2>/dev/null || true
sleep 1

exit 0
