# immich/rootfs/etc/s6-overlay/s6-rc.d/svc-server/finish
#!/usr/bin/with-contenv bash
# Kill orphaned Immich worker processes before s6 restarts the service

if pgrep -f "/usr/src/app" > /dev/null 2>&1; then
    echo "svc-server/finish: killing stale Immich worker processes" >&2
    pkill -9 -f "/usr/src/app" 2>/dev/null || true
    sleep 1
fi
