# immich/rootfs/etc/s6-overlay/s6-rc.d/svc-server/run
#!/usr/bin/with-contenv bash
set -euo pipefail

source /usr/local/bin/load-env.sh

exec gosu "${PUID}:${PGID}" /usr/local/bin/immich-server-wrapper.sh
