# karakeep/rootfs/etc/s6-overlay/s6-rc.d/init-config/run
#!/usr/bin/with-contenv bash
set -euo pipefail

source /usr/local/bin/load-env.sh

mkdir -p \
  /data \
  /data/meilisearch \
  /data/assets \
  /config
