# immich/rootfs/etc/s6-overlay/s6-rc.d/init-adduser/run
#!/usr/bin/with-contenv bash
set -euo pipefail

source /usr/local/bin/load-env.sh

if id immich >/dev/null 2>&1; then
  groupmod -o -g "${PGID}" immich || true
  usermod -o -u "${PUID}" -g "${PGID}" immich || true
fi

# Keep /data ownership intact for Postgres and Valkey runtime users.
chown -R "${PUID}:${PGID}" /config /defaults || true
