# karakeep/rootfs/etc/s6-overlay/s6-rc.d/init-db-migration/run
#!/usr/bin/with-contenv bash
set -euo pipefail

source /usr/local/bin/load-env.sh

echo "Running Karakeep database migration" >&2
cd /db_migrations
exec node index.js
