# librespot/rootfs/etc/services.d/librespot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs librespot
# ==============================================================================
declare -a options

if ! $(bashio::config.has_value 'name'); then
    bashio::log.fatal
    bashio::log.fatal "Add-on configuration is incomplete!"
    bashio::log.fatal
    bashio::log.fatal "The Librespot server needs to be identifiable with a name"
    bashio::log.fatal "and it seems you haven't configured one."
    bashio::log.fatal
    bashio::log.fatal "Please set the 'name' add-on option."
    bashio::log.fatal
    bashio::exit.nok
else
    options+=(--name $(bashio::config 'name'))
fi

if $(bashio::config.has_value 'bitrate'); then
    options+=(--bitrate $(bashio::config 'bitrate'))
fi
if $(bashio::config.has_value 'username') && $(bashio::config.has_value 'password'); then
    options+=(--username $(bashio::config 'username'))
    options+=(--password $(bashio::config 'password'))
fi
if $(bashio::config.has_value 'initial_volume'); then
    options+=(--initial-volume $(bashio::config 'initial_volume'))
fi
if $(bashio::config.true 'disable_audio_cache'); then
    options+=(--disable-audio-cache)
fi
if $(bashio::config.true 'linear_volume'); then
    options+=(--linear-volume)
fi
if $(bashio::config.true 'enable-volume-normalisation'); then
    options+=(--enable-volume-normalisation)
fi
if $(bashio::config.true 'disable_discovery'); then
    options+=(--disable-discovery)
fi
if $(bashio::config.true 'disable-gapless'); then
    options+=(--disable-gapless)
fi

# Are we running in debug mode?
if bashio::debug; then
    options+=(-v)
fi
# Pipe to snapserver
options+=(--backend pipe --device /share/snapfifo/librespot)

bashio::log.info 'Starting librespot...'

# Run librespot
exec librespot "${options[@]}"
