# qbittorrent/rootfs/etc/s6-overlay/s6-rc.d/init-qbittorrent/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: qBittorrent
# Configures qBittorrent
# ==============================================================================

readonly CONFIG_DIR="/config/qBittorrent/config"
readonly DATA_DIR="/config/qBittorrent/data"

bashio::log.info "Initializing qBittorrent..."

# Create configuration directory if it doesn't exist
if ! bashio::fs.directory_exists "${CONFIG_DIR}"; then
    bashio::log.debug "Creating config directory..."
    mkdir -p "${CONFIG_DIR}"
fi

# Create data directory if it doesn't exist
if ! bashio::fs.directory_exists "${DATA_DIR}"; then
    bashio::log.debug "Creating data directory..."
    mkdir -p "${DATA_DIR}"
fi

# Ensure directories are writable
chmod -R 755 /config/qBittorrent

bashio::log.info "Initialization complete!"
