# qbittorrent/rootfs/etc/s6-overlay/s6-rc.d/init-qbittorrent/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Configures qBittorrent
# ==============================================================================

readonly CONFIG_DIR="/data/qBittorrent"
readonly CONFIG_FILE="${CONFIG_DIR}/qBittorrent.conf"

bashio::log.info "Initializing qBittorrent..."

# Create configuration directory if it doesn't exist
if ! bashio::fs.directory_exists "${CONFIG_DIR}"; then
    bashio::log.debug "Creating config directory..."
    mkdir -p "${CONFIG_DIR}"
fi

# Create initial config file only if it doesn't exist
if [[ ! -f "${CONFIG_FILE}" ]]; then
    bashio::log.info "Creating initial qBittorrent configuration for reverse proxy access..."
    cat > "${CONFIG_FILE}" <<EOF
[Preferences]
WebUI\Address=*
WebUI\Port=8080
WebUI\LocalHostAuth=false
WebUI\CSRFProtection=false
WebUI\HostHeaderValidation=false
WebUI\ServerDomains=*
WebUI\UseUPnP=false

[LegalNotice]
Accepted=true
EOF
fi

bashio::log.info "Initialization complete!"
