# shairport-sync/rootfs/etc/services.d/shairport-sync/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs shairport-sync
# ==============================================================================
declare -a options

if ! $(bashio::config.has_value 'name'); then
    bashio::log.fatal
    bashio::log.fatal "Add-on configuration is incomplete!"
    bashio::log.fatal
    bashio::log.fatal "The AirPlay client needs to be identifiable with a name"
    bashio::log.fatal "and it seems you haven't configured one."
    bashio::log.fatal
    bashio::log.fatal "Please set the 'name' add-on option."
    bashio::log.fatal
    bashio::exit.nok
else
    options+=(-a $(bashio::config 'name'))
fi

# Are we running in debug mode?
if bashio::debug; then
  options+=(-vvv)
fi

# Log to stdout
options+=(-u)

# Pipe to snapserver
options+=(-o pipe -- /share/snapcast/shairport_sync)

bashio::log.info 'Starting shairport-sync...'

# Run shairport-sync
exec shairport-sync "${options[@]}"
