# unbound/rootfs/etc/services.d/unbound/run
#!/usr/bin/with-contenv bashio
CONFIG="/etc/unbound/unbound.conf"

# Run unbound
bashio::log.info "Starting unbound..."
exec unbound -c "${CONFIG}" -d
