# gluster/rootfs/etc/services.d/gluster/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs snapclient
# ==============================================================================
bashio::log.info 'Starting gluster...'

# Run snapclient
exec /usr/sbin/glusterd -N
