# immich/rootfs/etc/s6-overlay/s6-rc.d/immich-media-folder-migration/up
/etc/s6-overlay/s6-rc.d/immich-media-folder-migration/run
