# syncthing/root/etc/s6-overlay/s6-rc.d/syncthing-setup/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

bashio::net.wait_for 8384 127.0.0.1
bashio::log.info 'Post-start Syncthing configuration'

CONFIG_DIR="${STCONFDIR:-/data/config}"

if [ -f "$CONFIG_DIR/config.xml" ]; then
	apikey=$(sed -n 's:.*<apikey>\(.*\)</apikey>.*:\1:p' "$CONFIG_DIR/config.xml" | head -n1 || true)
	if [ -n "${apikey:-}" ]; then
		CLI_BASE="syncthing cli --gui-address=127.0.0.1:8384 --gui-apikey=$apikey"

		bashio::log.info 'Configuring Ingress compatibility settings'

		# Enable frame loading for Ingress
		$CLI_BASE config gui insecure-allow-frame-loading set true || \
			bashio::log.warning "Failed to set insecure-allow-frame-loading"

		# Skip hostname check for reverse proxy
		$CLI_BASE config gui insecure-skip-host-check set true || \
			bashio::log.warning "Failed to set insecure-skip-host-check"

		# Set default folder path
		bashio::log.info 'Setting default folder path to /data/sync'
		$CLI_BASE config defaults folder path set /data/sync || \
			bashio::log.warning "Failed to set default folder path"

		bashio::log.info 'Syncthing configuration completed'
	else
		bashio::log.warning 'API key not found; skipping setup'
	fi
else
	bashio::log.warning 'config.xml not found; skipping setup'
fi
