# murmur/rootfs/etc/services.d/murmur/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs Murmur
# ==============================================================================

bashio::log.info 'Starting Murmur...'

CONFIGFILE="/etc/murmur/murmur.ini"

# Run murmur 
exec /usr/bin/murmurd -fg -ini "${CONFIGFILE}"
