# snapserver/rootfs/etc/services.d/spotifyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs snapserver
# ==============================================================================
bashio::log.info 'Starting Snapserver...'

# Run snapserver
exec /usr/bin/snapserver
