# garage/rootfs/etc/s6-overlay/s6-rc.d/init-garage/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Configures Garage S3
# ==============================================================================

readonly CONFIG_PATH="/data/garage.toml"
readonly METADATA_DIR="/data/meta"
readonly DATA_DIR="/data/data"
readonly SECRETS_FILE="/data/.garage_secrets"

bashio::log.info "Initializing Garage S3..."

# Create directories
mkdir -p "${METADATA_DIR}" "${DATA_DIR}"

# Load or generate persistent secrets
load_or_generate_secret() {
    local key="$1"
    local value=""

    # Check if secrets file exists and contains this key
    if [ -f "${SECRETS_FILE}" ]; then
        value=$(grep "^${key}=" "${SECRETS_FILE}" 2>/dev/null | cut -d'=' -f2-) || value=""
    fi

    # If no stored value, generate new
    if [ -z "${value}" ]; then
        if [ "${key}" = "rpc_secret" ]; then
            value=$(openssl rand -hex 32)
        else
            value=$(openssl rand -base64 32)
        fi
        bashio::log.info "Generated new ${key}"

        # Persist to secrets file
        if [ -f "${SECRETS_FILE}" ]; then
            # Remove old entry if exists, then append
            grep -v "^${key}=" "${SECRETS_FILE}" > "${SECRETS_FILE}.tmp" 2>/dev/null || true
            mv "${SECRETS_FILE}.tmp" "${SECRETS_FILE}"
        fi
        echo "${key}=${value}" >> "${SECRETS_FILE}"
        chmod 600 "${SECRETS_FILE}"
    fi

    echo "${value}"
}

# Load or generate secrets (persistent across restarts)
RPC_SECRET=$(load_or_generate_secret "rpc_secret")
ADMIN_TOKEN=$(load_or_generate_secret "admin_token")
METRICS_TOKEN=$(load_or_generate_secret "metrics_token")

bashio::log.info "Using persistent secrets from ${SECRETS_FILE}"

# Get S3 region
S3_REGION=$(bashio::config 's3_region')

# Generate garage.toml
bashio::log.info "Generating Garage configuration..."
cat > "${CONFIG_PATH}" <<EOF
metadata_dir = "${METADATA_DIR}"
data_dir = "${DATA_DIR}"
db_engine = "sqlite"
replication_factor = 1

rpc_bind_addr = "[::]:3901"
rpc_public_addr = "127.0.0.1:3901"
rpc_secret = "${RPC_SECRET}"

[s3_api]
s3_region = "${S3_REGION}"
api_bind_addr = "[::]:3900"
root_domain = ".s3.garage"

[admin]
api_bind_addr = "[::]:3903"
admin_token = "${ADMIN_TOKEN}"
metrics_token = "${METRICS_TOKEN}"
EOF

# Start Garage server in background to initialize
bashio::log.info "Starting Garage server for initialization..."
garage -c "${CONFIG_PATH}" server &
GARAGE_PID=$!

# Wait for Garage to be ready
sleep 5

# File to store API keys persistently
readonly KEYS_FILE="/data/.garage_keys"

# Function to save key credentials
save_key_credentials() {
    local key_name="$1"
    local key_info

    key_info=$(garage -c "${CONFIG_PATH}" key info "${key_name}" 2>/dev/null)
    if [ -n "${key_info}" ]; then
        local key_id secret_key
        key_id=$(echo "${key_info}" | grep "Key ID:" | awk '{print $3}')
        secret_key=$(echo "${key_info}" | grep "Secret key:" | awk '{print $3}')

        if [ -n "${key_id}" ] && [ -n "${secret_key}" ]; then
            # Remove old entry for this key name if exists
            if [ -f "${KEYS_FILE}" ]; then
                grep -v "^${key_name}|" "${KEYS_FILE}" > "${KEYS_FILE}.tmp" 2>/dev/null || true
                mv "${KEYS_FILE}.tmp" "${KEYS_FILE}"
            fi
            # Save in format: name|key_id|secret_key
            echo "${key_name}|${key_id}|${secret_key}" >> "${KEYS_FILE}"
            chmod 600 "${KEYS_FILE}"
            bashio::log.info "Saved credentials for key: ${key_name}"
        fi
    fi
}

# Function to export keys to user-accessible location and log them
export_and_log_keys() {
    local export_file="/share/garage_keys.txt"

    if [ -f "${KEYS_FILE}" ]; then
        bashio::log.info "=== Garage S3 API Keys ==="

        # Write header to export file
        cat > "${export_file}" <<EOF
# Garage S3 API Keys
# Generated by Garage S3 addon
# Use these credentials to connect your S3 clients
#
EOF

        while IFS='|' read -r name key_id secret_key; do
            # Log to addon logs
            bashio::log.info "Key: ${name}"
            bashio::log.info "  Key ID: ${key_id}"
            bashio::log.info "  Secret: ${secret_key}"

            # Write to export file
            cat >> "${export_file}" <<EOF
[${name}]
Key ID: ${key_id}
Secret: ${secret_key}

EOF
        done < "${KEYS_FILE}"

        bashio::log.info "==========================================="
        bashio::log.info "Keys exported to: ${export_file}"
        bashio::log.info "Access via File Editor, VS Code, or Samba"
    fi
}

# Initialize cluster layout if needed
initialize_layout() {
    local node_id
    node_id=$(garage -c "${CONFIG_PATH}" status 2>/dev/null | grep "127.0.0.1:3901" | awk '{print $1}')

    if [ -z "${node_id}" ]; then
        bashio::log.warning "Could not get node ID, will retry on next start"
        return 1
    fi

    bashio::log.info "Node ID: ${node_id}"

    # Check if layout is already applied
    if garage -c "${CONFIG_PATH}" layout show 2>/dev/null | grep -q "${node_id}"; then
        bashio::log.info "Layout already configured"
        return 0
    fi

    # Assign and apply layout
    bashio::log.info "Assigning layout to node ${node_id}..."
    garage -c "${CONFIG_PATH}" layout assign -z dc1 -c 1G "${node_id}" || true

    bashio::log.info "Applying layout..."
    garage -c "${CONFIG_PATH}" layout apply --version 1 || true

    # Wait for layout to propagate
    bashio::log.info "Waiting for layout to propagate..."
    sleep 5

    return 0
}

# Ensure keys and buckets exist (runs every startup to handle config changes)
ensure_keys_and_buckets() {
    bashio::log.info "Ensuring keys and buckets from configuration..."
    local keys_len
    keys_len=$(bashio::config 'keys|length')

    for i in $(seq 0 $((keys_len-1))); do
        local name
        name=$(bashio::config "keys[${i}].name")

        bashio::log.info "Ensuring key exists: ${name}"

        # Check if key exists in Garage
        if ! garage -c "${CONFIG_PATH}" key info "${name}" >/dev/null 2>&1; then
            # Key doesn't exist - check if we have saved credentials
            if [ -f "${KEYS_FILE}" ]; then
                local saved_line
                saved_line=$(grep "^${name}|" "${KEYS_FILE}" 2>/dev/null || true)
                if [ -n "${saved_line}" ]; then
                    local saved_id saved_secret
                    saved_id=$(echo "${saved_line}" | cut -d'|' -f2)
                    saved_secret=$(echo "${saved_line}" | cut -d'|' -f3)
                    bashio::log.info "Recreating key ${name} with saved credentials..."
                    # Import the key with saved ID and secret
                    if ! garage -c "${CONFIG_PATH}" key import --yes "${saved_id}" "${saved_secret}" --name "${name}" 2>/dev/null; then
                        bashio::log.warning "Failed to import key, creating new: ${name}"
                        if ! garage -c "${CONFIG_PATH}" key create "${name}"; then
                            bashio::log.error "Failed to create key: ${name}"
                            continue
                        fi
                    fi
                else
                    bashio::log.info "Creating new key: ${name}"
                    garage -c "${CONFIG_PATH}" key create "${name}"
                fi
            else
                bashio::log.info "Creating new key: ${name}"
                if ! garage -c "${CONFIG_PATH}" key create "${name}"; then
                    bashio::log.error "Failed to create key: ${name}"
                    continue
                fi
            fi
            # Save credentials after creation
            save_key_credentials "${name}"
        else
            # Key exists - ensure we have it saved
            if [ ! -f "${KEYS_FILE}" ] || ! grep -q "^${name}|" "${KEYS_FILE}" 2>/dev/null; then
                save_key_credentials "${name}"
            fi
        fi

        # Ensure buckets and permissions
        local buckets_len
        buckets_len=$(bashio::config "keys[${i}].buckets|length")
        for j in $(seq 0 $((buckets_len-1))); do
            local bucket
            bucket=$(bashio::config "keys[${i}].buckets[${j}]")
            bashio::log.info "Ensuring bucket + allow: ${bucket} -> ${name}"

            garage -c "${CONFIG_PATH}" bucket create "${bucket}" >/dev/null 2>&1 || true
            garage -c "${CONFIG_PATH}" bucket allow \
                --read --write --owner \
                "${bucket}" \
                --key "${name}" 2>/dev/null || true
        done
    done
}

# Main initialization logic
bashio::log.info "Checking cluster state..."

if initialize_layout; then
    ensure_keys_and_buckets
    export_and_log_keys
    touch /data/.initialized
    bashio::log.info "Initialization complete!"
else
    bashio::log.warning "Layout initialization failed, will retry on next start"
fi

# Stop the background server
kill ${GARAGE_PID} 2>/dev/null || true
wait ${GARAGE_PID} 2>/dev/null || true

bashio::log.info "Initialization finished"
