# snapserver/rootfs/etc/services.d/snapserver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs snapserver
# ==============================================================================

## Run snapserver
bashio::log.info 'Starting Snapserver...'
exec /usr/bin/snapserver_up
