# gitea/rootfs/etc/s6-overlay/s6-rc.d/gitea/finish
#!/command/execlineb -S0
s6-svscanctl -t /run/service
