# gitea/rootfs/etc/s6-overlay/s6-rc.d/gitea/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Gitea..."

export GITEA_WORK_DIR=/data/gitea
export HOME=/data/gitea

cd /data/gitea
chown -R git:git /data/gitea

exec s6-setuidgid git gitea web --config /data/gitea/custom/conf/app.ini
