# snapclient/rootfs/etc/services.d/snapclient/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs snapclient
# ==============================================================================
bashio::log.info 'Starting snapclient...'

# Run snapclient
exec /usr/bin/snapclient
