# syncthing/root/etc/s6-overlay/s6-rc.d/syncthing/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info 'Intializing Syncthing'

ip=0.0.0.0

bashio::log.info 'Starting Syncthing'
syncthing --gui-address="$ip:8384" --no-browser
