# syncthing/root/etc/s6-overlay/s6-rc.d/syncthing/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

bashio::log.info 'Initializing Syncthing service'

CONFIG_DIR="${STCONFDIR:-/data/config}"
STATE_DIR="${STDATADIR:-/data/state}"
SYNC_ROOT="/data/sync"

mkdir -p "$CONFIG_DIR" "$STATE_DIR" "$SYNC_ROOT"
chmod 700 "$CONFIG_DIR" || true

# Get credentials from addon options
gui_user=$(bashio::config 'gui_user' || echo "")
gui_password=$(bashio::config 'gui_password' || echo "")

# Require both username and password
if [ -z "${gui_user:-}" ] || [ -z "${gui_password:-}" ]; then
	bashio::log.fatal "GUI username and password must be configured!"
	bashio::log.fatal "Please set 'gui_user' and 'gui_password' in addon configuration"
	bashio::exit.nok
fi

# Set or update credentials
bashio::log.info 'Configuring Syncthing with credentials'
echo "$gui_password" | syncthing generate \
	--config="$CONFIG_DIR" \
	--gui-user="$gui_user" \
	--gui-password=- || true

GUI_ADDR="0.0.0.0:8384"

bashio::log.info 'Starting Syncthing'
exec syncthing \
	--gui-address="$GUI_ADDR" \
	--no-browser
