# gluster/rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Starts the SSH service
# ============================================================================== 

# Start the SSH daemon
exec /usr/sbin/sshd < /dev/null
