# gluster/rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Starts the SSH service
# ============================================================================== 

bashio::log.info "Starting the SSH daemon..."
exec /usr/sbin/sshd -D -e
